# usr/share/upstart/sessions/yavdr-wm.conf
start on desktop-start
stop on desktop-end

respawn

exec /usr/bin/yavdr-wm
