# usr/share/upstart/sessions/yavdr-wm.conf
start on started dbus and started gnome-session
stop on stopping gnome-session

respawn

exec su - $USER -c /usr/bin/yavdr-wm
