# usr/share/upstart/sessions/yavdr-wm.conf
start on xsession
stop on session-end

respawn

exec /usr/bin/yavdr-wm
