# etc/init/vdr-post-start.conf
start on started vdr

task

script
  LANG=C /sbin/status lightdm | grep -q process && /sbin/initctl emit vdr-frontend-restart || /bin/true
end script
