# etc/init/yavdr-wm.conf
start on yavdr-desktop-session-start
stop on desktop-shutdown

respawn

exec su - $USER -c /usr/bin/yavdr-wm
